# node-agent/systemd-roomcast-agent.service
[Unit]
Description=RoomCast Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Environment=AGENT_SHARED_SECRET=changeme
Environment=MIXER_CONTROL=Master
WorkingDirectory=/home/pi/roomcast/node-agent
ExecStart=/usr/bin/python3 /home/pi/roomcast/node-agent/agent.py
Restart=always
RestartSec=3
User=pi

[Install]
WantedBy=multi-user.target
